# config/systemd/corporate-assistant.service
[Unit]
Description=Corporate Assistant Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=assistant
Group=assistant
WorkingDirectory=/opt/corporate-assistant
EnvironmentFile=/opt/corporate-assistant/.env
ExecStart=/opt/corporate-assistant/.venv/bin/assistant-bot bot
ExecStartPre=/opt/corporate-assistant/.venv/bin/assistant-bot doctor --quiet
Restart=always
RestartSec=5
TimeoutStopSec=30
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/corporate-assistant/data /opt/corporate-assistant/logs /opt/corporate-assistant/backups

[Install]
WantedBy=multi-user.target
